# system/vpn.socket
[Unit]
Description=Virbox VPN WSGI socket
[Socket]
ListenStream=/run/vpn.sock
[Install]
WantedBy=sockets.target

